# bash_scripts/pc_info_log_service/pc_info_service.timer
[Unit]
Description=timer for system monitoring

[Timer]
OnUnitActiveSec=30s
OnBootSec=10s

[Install]
WantedBy=timers.target
